# debian/cron.d
@reboot		/usr/bin/screen-auto-start
*/5 * * * *	/usr/bin/screen-save-all-sessions
