# debian/screen-save.screen-save-autostart.user.service
[Unit]
Description=Auto start all saved screen session

[Service]
ExecStart=/usr/bin/screen-auto-start

[Install]
WantedBy=default.target