# debian/screen-save.screen-save-all.timer
[Unit]
Description=Save all 5' minutes all screen session

[Timer]
OnCalendar=*:0/15
OnStartupSec=15min

[Install]
WantedBy=timers.target
