# debian/screen-save.screen-save-all.service
[Unit]
Description=Save all screen session

[Service]
Type=oneshot
ExecStart=/usr/bin/screen-save-all-sessions
User=root
