# debian/screen-save.screen-save-autostart.timer
[Unit]
Description=Auto start all saved session

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
