# debian/screen-auto-start.service
[Unit]
Description=Auto start all saved screen session

[Service]
Type=oneshot
ExecStart=/usr/bin/screen-auto-start
User=root
